# server/.crontab
# Keep dotfiles in sync
*/1 * * * * cd /root/Documents/code/dotfiles && git pull 

# Run certbot renewal daily at 3, since it only fires if certs are close to expiry
0 3 */1 * * /usr/bing/certbot renew --pre-hook "service nginx stop" --post-hook "service nginx start" --quiet --no-self-upgrade 
