# server/crontab
# Refresh startpage repo
*/1 * * * * cd /var/www/startpage && git pull

# Keep dotfiles in sync
*/1 * * * * cd /root/Documents/code/dotfiles && git pull 
