# server/init/actions_runner.service
[Unit]
Description=actions_runner service

[Service]
Type=simple
User=chl
Restart=always
ExecStart=/home/chl/Documents/code/actions-runner/run.sh
